# creating_custom_deb_package/helloworld-1.0-1/DEBIAN/postinst
#!/bin/bash

systemctl enable helloworld-daemon.service
systemctl start helloworld-daemon.service
