# creating_custom_deb_package/helloworld-1.0-1/DEBIAN/prerm
#!/bin/bash

systemctl stop helloworld-daemon.service
systemctl disable helloworld-daemon.service
